# services/accuweather.timer
[Unit]
Description=Run accuweather logger service periodically

[Timer]
Unit=accuweather.service
OnCalendar=*:0/30

[Install]
WantedBy=timers.target